# flannel/scripts/deb/systemd/after-install.sh
        # Initial installation
        systemctl preset flanneld.service >/dev/null 2>&1 || :
